# build/config/supervisor/crond.conf
[program:crond]
autostart=true
autorestart=true
redirect_stderr=true
command=/usr/sbin/crond -f -c /var/spool/cron/crontabs
numprocs=1
