# php73/supervisord/sshd.conf
[program:sshd]
autostart=true
autorestart=true
redirect_stderr=true
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config
numprocs=1